# extras/autostart/kiosk-epiphany.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=Start kiosk
Encoding=UTF-8
Exec=sh -c "/home/pi/kiosk-browser.sh"
X-GNOME-Autostart-enabled=true
